# distfiles/mtos/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXH9XiwF installed on Sun May  6 22:46:50 2012)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
0,20,40 * * * * cd /var/www/localhost/cgi-bin/mt && echo "./tools/run-periodic-tasks" | su --shell=/bin/sh apache

