# distfiles/openpne/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXH9XiwF installed on Sun May  6 22:46:50 2012)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
00   6 * * * sh /var/www/localhost/htdocs/bin/send_daily_news.cron /var/www/localhost/htdocs /usr/bin/php
00   6 * * * sh /var/www/localhost/htdocs/bin/birthday_mail.cron /var/www/localhost/htdocs /usr/bin/php

