# distfiles/redmine/generate-secret-token
#!/sbin/runscript

depend() {
	need mysql
}

start () {
	ebegin "Generating secret token for redmine..."
	if [ ! -f /var/www/localhost/htdocs/config/initializers/session_store.rb ]; then
		cd /var/www/localhost/htdocs && rake db:migrate RAILS_ENV=production
		cd /var/www/localhost/htdocs && rake generate_secret_token && chown apache.apache /var/www/localhost/htdocs/config/initializers/secret_token.rb
		mysql redmine -e "update users set language='ja' where id=1"
		mysql redmine -e "insert into settings(name,value,updated_on) values('default_language','ja',now())"
	fi
	eend $?
}

