# distfiles/redmine/generate-secret-token
#!/sbin/runscript

start () {
	ebegin "Generating secret token for redmine..."
	if [ ! -f /var/www/localhost/htdocs/config/initializers/session_store.rb ]; then
		cd /var/www/localhost/htdocs && rake generate_secret_token && chown apache.apache /var/www/localhost/htdocs/config/initializers/secret_token.rb
	fi
	eend $?
}

