# distfiles/sugarcrm/crontab
* * * * * cd /var/www/localhost/htdocs; php -f cron.php > /dev/null 2>&1 
